# files/ecr.ini
[supervisord]
nodaemon=true
pidfile=/supervisor/run/supervisord.pid
logfile=/supervisor/logs/supervisord.log
childlogdir=/supervisor/logs

[program:nginx]
command = nginx -g "daemon off;"
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
autostart = true

[program:crond]
command = /scripts/setup-cron.sh
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
autostart = true
